# src/etc/NetworkManager/dispatcher.d/00-hostname
#!/bin/sh

iface=$1
action=$2

case ${action} in
hostname)
        logger -p 4 -t nm.dispatch "Hostname changed on ${iface}"
        ;;
esac
