# src/etc/NetworkManager/dispatcher.d/10-openrc-status
# The OpenRC status script provided by Gentoo is disabled,
# since it is a rather heavy hack, which may also interfere
# with, e.g., mac-randomize.
# (#252137)
