# src/etc/init.d/cabled
#!/sbin/runscript

daemon=/usr/local/libexec/cable/cabled
pidfile=/var/run/cabled.pid
user=anon

depend() {
    after  liberte
    before nginx
}

start() {
    ebegin "Starting cabled"
    start-stop-daemon -S -q -e CHUSER=${user} -p ${pidfile} \
                      -a ${daemon} -- fork
    eend $?
}

stop() {
    ebegin "Stopping cabled"
    start-stop-daemon -K -q -p ${pidfile}
    eend $?
}
