# src/etc/init.d/cabled
#!/sbin/runscript

cableflag=/var/www/cable
daemon=/usr/local/libexec/cable/cabled
pidfile=/var/run/cabled.pid
luser=anon

depend() {
    # need identity setup
    need identity

    # initialization in cabled is done before nginx is up
    before nginx
}

start() {
    if [ ! -e ${cableflag} ]; then
        ebegin "Starting cabled"
        start-stop-daemon -S -q -e CHUSER=${luser} -p ${pidfile} \
                          -a ${daemon} -- fork
        eend $?
    else
        ewarn "Skipping cabled startup"
    fi
}

stop() {
    if [ -e ${pidfile}  -o  ! -e ${cableflag} ]; then
        ebegin "Stopping cabled"
        start-stop-daemon -K -q -p ${pidfile}
        eend $?
    fi
}
