# src/etc/init.d/cabled
#!/sbin/runscript

description="Starts cables communication daemon."

luser=anon
cableflag=/var/www/cable

daemon=/usr/local/libexec/cable/cabled
pidfile=/var/run/cabled.pid

depend() {
    # need identity setup
    need identity

    # initialization in cabled is done before nginx is up
    before nginx
}

start() {
    if [ ! -e ${cableflag} ]; then
        # service_stopped doesn't work, since nginx is "starting"
        if ! service_started nginx; then
            ebegin "Initializing cabled"
            sudo -n -u ${luser} ${daemon} init
            eend $?
        else
            ewarn "Skipping cabled init (nginx not inactive)"
        fi

        # SSD -u sets up USER and HOME, and sanitized environment
        ebegin "Starting cabled"
        start-stop-daemon -S -q -u ${luser} -p ${pidfile} -bm -x ${daemon} --
        eend $?
    else
        ewarn "Skipping cabled startup"
    fi
}

stop() {
    if [ -e ${pidfile}  -o  ! -e ${cableflag} ]; then
        ebegin "Stopping cabled"
        start-stop-daemon -K -q -u ${luser} -p ${pidfile}
        eend $?
    fi
}
