# src/etc/init.d/lockdown
#!/sbin/runscript

description="Arms boot media removal watchdog."

boot=/mnt/boot

launcher=/usr/local/sbin/poweroff-watchdog
pidfile=/var/run/poweroff-watchdog.pid


depend() {
    # this is a boot-level service
    need localmount

    # using /var/run
    after bootmisc

    # fw-reload needs /var/run/usage-profile
    before iptables
}

start() {
    # Lock root account's password, unless gentoo=root was given
    if get_bootparam 'root'; then
        ebegin Enabling root password: liberte
        usermod -U root
        eend $?
    fi


    # Arm KEXEC wipe & halt
    ebegin Arming anti-forensic RAM wipe with halt
    kexec-load arm
    eend $?


    # Arm the poweroff watchdog
    if mountpoint -q ${boot}; then
        ebegin Arming power-off on boot media removal
        mediadev=/dev/block/`mountpoint -d ${boot}`
        start-stop-daemon -S -q -p ${pidfile} -bm -x ${launcher} -- ${mediadev}
        eend $?
    else
        ewarn Skipping boot media removal watchdog
    fi


    # Unblock [wifi wimax wwan], block [bluetooth uwb gps fm]
    if [ -e /dev/rfkill ]; then
        ebegin Blocking all radios except Wi-Fi
        rfkill block all
        rfkill unblock wifi && rfkill unblock wwan && rfkill unblock wimax
        eend $?
    fi


    # Support non-anonymous mode
    if get_bootparam 'noanon'; then
        ewarn Non-anonymous profile has been selected via boot options
        ebegin Initializing non-anonymous profile
        init-profile 'noanon'
        eend $?
    fi
}

stop() {
    # Disarm the poweroff watchdog
    if [ -e ${pidfile} ]; then
        ebegin Disarming power-off on boot media removal
        start-stop-daemon -K -q -p ${pidfile}
        eend $?
    fi
}
