# src/etc/init.d/lockdown
#!/sbin/runscript

description="Arms boot media removal watchdog."

boot=/mnt/boot

launcher=/usr/local/sbin/poweroff-watchdog
pidfile=/var/run/poweroff-watchdog.pid


depend() {
    # this is a boot-level service
    need localmount

    # using /var/run
    after bootmisc

    # fw-reload needs /var/run/usage-profile,
    # and also extracts IPs from /etc/hosts
    before iptables
}

start() {
    # Lock root account's password, unless gentoo=root was given
    if get_bootparam 'root'; then
        ebegin Enabling root password: liberte
        usermod -U root
        eend $?
    fi


    # Arm KEXEC wipe & halt
    ebegin Arming anti-forensic RAM wipe with halt
    kexec-load arm
    eend $?


    # Arm the poweroff watchdog
    ebegin Arming power-off on boot media removal
    mediadev=/dev/block/`mountpoint -d ${boot}`
    start-stop-daemon -S -p ${pidfile} -bm -x ${launcher} -- ${mediadev}
    eend $?


    # Unblock [wifi wwan], block [bluetooth wimax uwb gps fm]
    if [ -e /dev/rfkill ]; then
        ebegin Blocking all radios except Wi-Fi
        rfkill block all
        rfkill unblock wifi && rfkill unblock wwan
        eend $?
    fi


    # Disable NTP if requested
    if get_bootparam 'nontp'; then
        ewarn Disabling NTP time synchronization
        chmod a-x /etc/NetworkManager/dispatcher.d/01-ntp
        eend $?
    fi


    # Support non-anonymous mode
    if get_bootparam 'noanon'; then
        ewarn Non-anonymous profile has been selected via boot options
        ebegin Initializing non-anonymous profile
        init-profile 'noanon'
        eend $?
    fi
}

stop() {
    # Disarm the poweroff watchdog
    if [ -e ${pidfile} ]; then
        ebegin Disarming power-off on boot media removal
        start-stop-daemon -K -p ${pidfile}
        eend $?
    fi
}
