# src/etc/init.d/lockdown
#!/sbin/runscript

description="Arms boot media removal watchdog."

boot=/mnt/boot

launcher=/usr/local/sbin/poweroff-watchdog
pidfile=/var/run/poweroff-watchdog.pid

depend() {
    # power-off watchdog needs kexec
    use kexec

    # no point in shutting down the power-off watchdog
    keyword -stop
}

start() {
    # No mountpoint if there is no loop (direct mount on /mnt/live)
    cddevpath=$(udevadm info -q path -p /dev/block/$(mountpoint -d ${boot}) 2>/dev/null)


    # Lock root account's password, unless gentoo=root was given
    if ! get_bootparam 'root'; then
        ebegin Disabling root password
        usermod -L root
        eend $?
    fi


    # Arm the poweroff watchdog
    if [ -n "${cddevpath}" ]; then
        ebegin Arming power-off on boot media removal
        start-stop-daemon -S -q -p ${pidfile} -bm -x ${launcher} -- "${cddevpath}"
        eend $?
    else
        ewarn Skipping boot media removal watchdog
    fi
}

stop() {
    # Disarm the poweroff watchdog
    if [ -e ${pidfile} ]; then
        ebegin Disarming power-off on boot media removal
        start-stop-daemon -K -q -p ${pidfile}
        eend $?
    fi
}
