# src/etc/init.d/lockdown
#!/sbin/runscript

description="Arms boot media removal watchdog."

boot=/mnt/boot

cablepub=/srv/www
lidevents=/var/run/lid-events

launcher=/usr/local/sbin/poweroff-watchdog
pidfile=/var/run/poweroff-watchdog.pid


depend() {
    # this is a boot-level service
    need localmount

    # using /var/run
    after bootmisc

    # power-off watchdog uses kexec
    use kexec
}

start() {
    # Lock root account's password, unless gentoo=root was given
    if ! get_bootparam 'root'; then
        ebegin Disabling root password
        usermod -L root
        eend $?
    fi


    # Arm KEXEC wipe & halt
    ebegin Arming anti-forensic RAM wipe on shutdown
    kexec-load arm
    eend $?


    # Arm the poweroff watchdog
    if mountpoint -q ${boot}; then
        ebegin Arming power-off on boot media removal
        mediadev=/dev/block/`mountpoint -d ${boot}`
        start-stop-daemon -S -q -p ${pidfile} -bm -x ${launcher} -- ${mediadev}
        eend $?
    else
        ewarn Skipping boot media removal watchdog
    fi


    # See also: /etc/acpi/default.sh
    if [ ! -e ${lidevents} ]; then
        ebegin Creating closed lid indicator
        checkpath -q -d -m 750 -o root:video ${lidevents}
        touch ${lidevents}/close.flag
        eend $?
    fi


    # Workaround /run incompatibility (#400899)
    checkpath -q -d -m 700 /var/run/spawn-fcgi
}

stop() {
    # Disarm the poweroff watchdog
    if [ -e ${pidfile} ]; then
        ebegin Disarming power-off on boot media removal
        start-stop-daemon -K -q -p ${pidfile}
        eend $?
    fi
}
