# src/etc/init.d/lockdown
#!/sbin/runscript

description="Arms boot media removal watchdog and suspends disk drives."

boot=/mnt/boot

launcher=/usr/local/sbin/poweroff-watchdog
pidfile=/var/run/poweroff-watchdog.pid

depend() {
    # power-off watchdog needs kexec
    use kexec

    # no point in shutting down the power-off watchdog
    keyword -stop
}

start() {
    # No mountpoint if there is no loop (direct mount on /mnt/live)
    cddevpath=$(udevadm info -q path -p /dev/block/$(mountpoint -d ${boot}) 2>/dev/null)


    # Lock root account's password
    ebegin Disabling root password
    usermod -L root
    eend $?


    # Arm the poweroff watchdog
    if [ -n "${cddevpath}" ]; then
        ebegin Arming power-off on boot media removal
        start-stop-daemon -S -q -p ${pidfile} -bm -x ${launcher} -- "${cddevpath}"
        eend $?
    else
        ewarn Skipping boot media removal watchdog
    fi


    # Put ATA/SCSI drives into standby
    # /dev/disk/by-path should contain only physical devices
    if [ -e /dev/disk/by-path ]; then
        hdevs=`find /dev/disk/by-path ! -name '*-part[0-9]*' -xtype b  \
                                      -exec udevadm info -q path -n {} \;`
        for devpath in ${hdevs}; do
            dev=/dev/`udevadm info -q name -p ${devpath}`

            if hdparm -M 128 ${dev} 1>/dev/null 2>&1; then
                ebegin Enabled acoustic settings: ${dev}
                eend 0
            fi

            # Ignore boot media (desirable if it is a disk)
            if [ "${cddevpath#${devpath}}" = "${cddevpath}" ]; then
                if hdparm -y ${dev} 1>/dev/null 2>&1; then
                    ebegin Enabled standby mode: ${dev}
                    eend 0
                fi
            fi
        done
    fi
}

stop() {
    # Disarm the poweroff watchdog
    if [ -e ${pidfile} ]; then
        ebegin Disarming power-off on boot media removal
        start-stop-daemon -K -q -p ${pidfile}
        eend $?
    fi
}
