# src/etc/init.d/persist
#!/sbin/runscript

description="Restores and saves user configuration."

luser=anon
lmount=/home/${luser}/persist
persist=persistence
profile=

depend() {
    # need OTFE mount
    need liberte
}

start() {
    if mountpoint -q ${lmount}; then
        # Extract configuration archive, if present
        if get_bootparam 'nosettings'; then
            ewarn Skipping persistent configuration extraction
        else
            ebegin Restoring persistent configuration: ${lmount}/settings
            ! get_bootparam 'noanon' || profile=noanon
            sudo -n -u ${luser} ${persist} restore ${lmount}/settings ${profile}
            eend $?
        fi
    else
        ewarn Skipping persistent configuration extraction
    fi
}

stop() {
    if mountpoint -q ${lmount}; then
        # Save configuration (an incremental snapshot)
        if get_bootparam 'nosettings'; then
            ewarn Skipping persistent configuration saving
        else
            # ${lmount}/settings is created in identity service for the first time
            ebegin Saving configuration: ${lmount}/settings
            ! get_bootparam 'noanon' || profile=noanon
            sudo -n -u ${luser} ${persist} save ${lmount}/settings ${profile}
            eend $?
        fi
    else
        ewarn Skipping persistent configuration saving
    fi
}
