# src/etc/init.d/spindown
#!/sbin/runscript

description="Spins down and quiets down hard drives."

boot=/mnt/boot

depend() {
    # start after drive-probing services (autofs just in case)
    after laptop_mode smartd autofs

    # liberte needs user input, which suffices for drive activity to stop
    after liberte
}

start() {
    mediadev=$(readlink -f /dev/block/$(mountpoint -d ${boot}))

    # Put ATA/SCSI drives into quiet mode
    # /dev/disk/by-path does not cover all physical devices (some lack ID_PATH)
    # NOTE: hdparm -M and -y do not wake up disks from standby
    if [ -e /dev/disk/by-id ]; then
        ebegin Enabling drive acoustic and standby modes
        eindent

        devs=`find /dev/disk/by-id ! -name '*-part[0-9]*' -xtype b \
                                   -exec readlink -f {} \; | sort -u`
        for dev in ${devs}; do
            if hdparm -q -M 128 "${dev}" 2>/dev/null; then
                einfo "Quiet mode:   ${dev}"
            fi

            # Ignore boot media (desirable if it is a disk)
            if [ "${mediadev#${dev}}" = "${mediadev}" ]; then
                if hdparm -q -y "${dev}" 2>/dev/null; then
                    einfo "Standby mode: ${dev}"
                fi
            fi
        done

        eoutdent
        eend $?
    fi
}
