# src/etc/init.d/spindown
#!/sbin/runscript

description="Spins down and quiets down hard drives."

boot=/mnt/boot

depend() {
    # start after drive-probing services (autofs just in case)
    after laptop_mode smartd autofs

    # liberte needs user input, which suffices for drive activity to stop
    after liberte
}

start() {
    mediadev=none
    if mountpoint -q ${boot}; then
        mediadev=$(readlink -f /dev/block/$(mountpoint -d ${boot}))
    fi

    # Put ATA/SCSI drives into quiet mode
    # /dev/disk/by-path should contain only physical devices
    if [ -e /dev/disk/by-path ]; then
        devs=`find /dev/disk/by-path ! -name '*-part[0-9]*' -xtype b \
                                     -exec readlink -f {} \;`
        for dev in ${devs}; do
            if hdparm -M 128 ${dev} 1>/dev/null 2>&1; then
                ebegin Enabled acoustic settings: ${dev}
                eend 0
            fi

            # Ignore boot media (desirable if it is a disk)
            if [ "${mediadev#${dev}}" = "${mediadev}" ]; then
                if hdparm -y ${dev} 1>/dev/null 2>&1; then
                    ebegin Enabled standby mode: ${dev}
                    eend 0
                fi
            fi
        done
    fi
}
