# src/etc/init.d/tordate
#!/sbin/runscript

description="Sets date from Tor consensus and handles Tor/I2P restarts."

daemon=/usr/local/sbin/tordate
pidfile=/var/run/tordate.pid

depend() {
    before NetworkManager tor i2p
}

start() {
    ebegin "Starting tordate"
    start-stop-daemon -S -q -p ${pidfile} -bm -x ${daemon} --
    eend $?
}

stop() {
    ebegin "Stopping tordate"
    start-stop-daemon -K -q -p ${pidfile}
    eend $?
}
