# src/etc/init.d/tordate
#!/sbin/runscript

description="Sets date from Tor consensus and handles Tor restarts."

command=/usr/local/sbin/tordate
command_background="true"
pidfile=/var/run/tordate.pid

depend() {
    # need /var
    need   localmount

    before NetworkManager tor
}
