# src/etc/init.d/vmconfig
#!/sbin/runscript

description="Configures the system as a guest in virtual machine."

vmflag=/var/run/vmtype

cleanlabel() {
    local label="$1"
    echo "${label}" | tr '\000-\011\013-\047' _
}


depend() {
    # need /etc
    need localmount

    # using /var/run
    after bootmisc

    # don't bother autofs with SIGHUPs from ps-mount
    before autofs

    # user's X wallpaper depends on vmflag availability
    before xdm

    # this is a one-time setup (also, removal won't work correctly with spaces/quotes in DEVNAME)
    keyword -stop
}

start() {
    # VM flag is used for X wallpaper
    if [ ! -e ${vmflag} ]; then
        ebegin Detecting virtualization state
        virt-what > ${vmflag}
        eend $?
    fi


    if vmware-hgfsclient 1>/dev/null 2>&1; then
        ebegin Adding VMware shares to media automounts
        eindent

        # ID_FS_LABEL is quoted when in ps-mount
        vmware-hgfsclient | while read vmshare; do
            einfo /media/VMware-"${vmshare}"
            ACTION=add ID_FS_TYPE=vmhgfs DEVNAME=".host:/\"${vmshare}\"" ID_FS_LABEL=VMware-"$(cleanlabel "${vmshare}")" ps-mount
        done

        eoutdent
        eend $?
    fi


    if VBoxControl -nologo guestproperty get /VirtualBox/HostInfo/VBoxVer 1>/dev/null 2>&1; then
        ebegin Adding VirtualBox shares to media automounts
        eindent

        # ID_FS_LABEL is quoted when in ps-mount
        VBoxControl -nologo sharedfolder list | sed -n 's/^[[:digit:]]\+ - //p' | while read vbshare; do
            einfo /media/VBox-"${vbshare}"
            ACTION=add ID_FS_TYPE=vboxsf DEVNAME="\"${vbshare}\"" ID_FS_LABEL=VBox-"$(cleanlabel "${vbshare}")" ps-mount
        done

        eoutdent
        eend $?
    fi
}
