# src/etc/init.d/xconfig
#!/sbin/runscript

description="Sets up cgroups, configures X server and sound."

cgroupmnt=/sys/fs/cgroup
cgrouprel=/usr/local/sbin/cgroup-release

vmflag=/etc/vmtype
silentflag=/etc/splash/silent/active

xorgrc=/etc/X11/xorg.conf

amixerrst=/usr/local/bin/reset-mixer
asoundrc=/etc/asound.conf

depend() {
    need   localmount
    before xdm
}

start() {
    # Setup cgroups
    if [ ! -e ${cgroupmnt}/user ]; then
        ebegin Setting up control groups
        echo ${cgrouprel} > ${cgroupmnt}/release_agent
        echo 1            > ${cgroupmnt}/notify_on_release
        mkdir -m 3770 ${cgroupmnt}/user
        chgrp users   ${cgroupmnt}/user
        eend $?
    fi

    # VM flag is used for X wallpaper
    if [ ! -e ${vmflag} ]; then
        ebegin Detecting virtualization state
        virt-what > ${vmflag}
        eend $?
    fi

    if egrep -q '\<splash=([^ ]*,)?theme:silent\>' /proc/cmdline; then
        ebegin Detected silent theme
        touch ${silentflag}
        eend $?
    fi

    # Select preferred audio card based on the number of mixer controls
    if [ -e /proc/asound ]; then
        ebegin Selecting preferred audio card
        acards=`sed -n 's/^[[:blank:]]*[[:digit:]]\+[[:blank:]]*\[\([^][:blank:]]\+\)[[:blank:]]*\].*$/\1/p' /proc/asound/cards`
        bestcount=-1
        for acard in ${acards}; do
            acount=`amixer -c ${acard} controls | wc -l`
            if [ ${acount} -gt ${bestcount} ]; then
                bestcard=${acard}
                bestcount=${acount}
            fi
        done
        if [ ${bestcount} -gt -1 ]; then
            sed -i "s/\<card .*\$/card ${bestcard}/" ${asoundrc}
        else
            false
        fi
        eend $?

        ebegin Adjusting audio mixer values
        ${amixerrst}
        eend $?
    fi

    # Configure X server, unless "nox" is given to kernel
    if get_bootparam 'nox'; then
        ewarn Skipping X server configuration
    elif [ ! -e ${xorgrc} ]; then
        ebegin Configuring X server
        if HOME=/root Xorg -configure 2>/dev/null; then
            # Work around Xorg autodetection / fallback bug
            if grep -q 'Number of created screens does not match number of detected devices' /var/log/Xorg.0.log \
                && grep -q '"vesa"' /root/xorg.conf.new; then
                rm /root/xorg.conf.new /usr/lib/xorg/modules/drivers/vesa_drv.so
                HOME=/root Xorg -configure 2>/dev/null
            fi

            if [ -e /root/xorg.conf.new ]; then
                mv /root/xorg.conf.new ${xorgrc}

                # Bug #3319525 workaround (doesn't always work)
                touch /mnt/hidden/rwroot/etc/X11

                # Turn off text cursor on vt3 (cf. xserverrc)
                echo -n "\033[?25l" > /dev/tty3
            else
                false
            fi
        else
            false
        fi
        eend $?
    fi
}
