# src/etc/init.d/xconfig
#!/sbin/runscript

description="Configures X server and sound."

silentflag=/var/run/nologo

xorgdir=/etc/X11
xorglog=/var/log/Xorg.0.log
xorgdrv=/usr/lib/xorg/modules/drivers

amixerrst=/usr/local/bin/reset-mixer
asoundrc=/etc/asound.conf


depend() {
    # need /etc mount
    need   localmount

    # using /var/run
    after bootmisc

    # X server is configured here
    before xdm
}


start() {
    if get_bootparam 'nologo'; then
        ebegin Disabling desktop branding
        touch ${silentflag}
        eend $?
    fi


    # Select preferred audio card based on the number of mixer controls
    if [ -e /proc/asound ]; then
        ebegin Selecting preferred audio card
        acards=`sed -n 's/^[[:blank:]]*[[:digit:]]\+[[:blank:]]*\[\([^][:blank:]]\+\)[[:blank:]]*\].*$/\1/p' /proc/asound/cards`
        bestcount=-1
        for acard in ${acards}; do
            acount=`amixer -c ${acard} controls | wc -l`
            if [ ${acount} -gt ${bestcount} ]; then
                bestcard=${acard}
                bestcount=${acount}
            fi
        done
        if [ ${bestcount} -gt -1 ]; then
            sed -i "s/\<card .*\$/card ${bestcard}/" ${asoundrc}
        else
            false
        fi
        eend $?

        ebegin Adjusting audio mixer values
        ${amixerrst}
        eend $?
    fi


    # Configure X server, unless "nox" is given to kernel
    if get_bootparam 'nox'; then
        ewarn Skipping X server configuration
    elif [ ! -e ${xorgdir}/xorg.conf ]; then
        # Temporarily mount a VESA-only drivers directory if requested
        if get_bootparam 'xvesa'; then
            ewarn Forcing VESA X server driver
            mount ${xorgdrv}
        fi

        ebegin Configuring X server
        if HOME=${xorgdir} Xorg -configure 2>/dev/null; then
            # Work around Xorg autodetection / fallback bug
            if grep -q 'Number of created screens does not match number of detected devices' ${xorglog} \
                && grep -q '"vesa"' ${xorgdir}/xorg.conf.new && [ -e ${xorgdrv}/dummy_drv.so ]; then
                # Temporarily blacklist VESA driver and rerun autodetection
                mount  ${xorgdrv}/vesa_drv.so
                rm ${xorgdir}/xorg.conf.new
                HOME=${xorgdir} Xorg -configure 2>/dev/null
                umount ${xorgdrv}/vesa_drv.so
            fi

            if [ -e ${xorgdir}/xorg.conf.new ]; then
                mv ${xorgdir}/xorg.conf.new ${xorgdir}/xorg.conf

                # Turn off text cursor on vt3 (cf. xserverrc)
                echo -n "\033[?25l" > /dev/tty3
            else
                false
            fi
        else
            false
        fi
        eend $?

        # mountpoint fails for bind mounts
        if get_bootparam 'xvesa'; then
            umount ${xorgdrv}
        fi
    fi
}
