# src/etc/init.d/xconfig
#!/sbin/runscript

description="Configures X server and sound."

silentflag=/var/run/nologo

xorgdir=/etc/X11
xorglog=/var/log/Xorg.0.log
xorgdrv=/usr/lib/xorg/modules/drivers
xorgdrvtype=native

amixerrst=/usr/local/bin/reset-mixer
asoundrc=/etc/asound.conf


depend() {
    # need /etc mount
    need   localmount

    # using /var/run
    after bootmisc

    # X server is configured here
    before xdm
}


start() {
    if get_bootparam 'nologo'; then
        ebegin Disabling desktop branding
        touch ${silentflag}
        eend $?
    fi


    # Select preferred audio card based on the number of mixer controls
    if [ -e /proc/asound ]; then
        ebegin Selecting preferred audio card
        acards=`sed -n 's/^[[:blank:]]*[[:digit:]]\+[[:blank:]]*\[\([^][:blank:]]\+\)[[:blank:]]*\].*$/\1/p' /proc/asound/cards`
        bestcount=-1
        for acard in ${acards}; do
            acount=`amixer -c ${acard} controls | wc -l`
            if [ ${acount} -gt ${bestcount} ]; then
                bestcard=${acard}
                bestcount=${acount}
            fi
        done
        if [ ${bestcount} -gt -1 ]; then
            sed -i "s/\<card .*\$/card ${bestcard}/" ${asoundrc}
        else
            false
        fi
        eend $?

        ebegin Adjusting audio mixer values
        ${amixerrst}
        eend $?
    fi


    # Configure X server, unless "nox" is given to kernel
    if get_bootparam 'nox'; then
        ewarn Skipping X server configuration
    else
        # Temporarily mount native/VESA/fbdev-only drivers directory if requested
        if get_bootparam 'xvesa'; then
            ewarn Forcing VESA X server driver
            xorgdrvtype=vesa
        elif get_bootparam 'xfb'; then
            ewarn Forcing framebuffer X server driver
            xorgdrvtype=fbdev
        fi
        mount -rB ${xorgdrv}.${xorgdrvtype} ${xorgdrv}

        ebegin Configuring X server
        rm -f ${xorgdir}/xorg.conf.new
        HOME=${xorgdir} Xorg -configure 2>/dev/null
        if [ ! -e ${xorgdir}/xorg.conf.new  -a  ${xorgdrvtype} = native ]; then
            ewarn Falling back to VESA X server driver
            umount ${xorgdrv} && mount -rB ${xorgdrv}.vesa ${xorgdrv}
            HOME=${xorgdir} Xorg -configure 2>/dev/null
        fi
        umount ${xorgdrv}

        if [ -e ${xorgdir}/xorg.conf.new ]  &&  ! grep -q 'Configuration failed' ${xorglog}; then
            mv ${xorgdir}/xorg.conf.new ${xorgdir}/xorg.conf

            # Turn off text cursor on vt3 (cf. xserverrc)
            echo ${BASH:+-e} -n "\033[?25l" > /dev/tty3
        else
            false
        fi
        eend $?
    fi
}
