# src/etc/init.d/xconfig
#!/sbin/runscript

description="Configures X server and sound."

silentflag=/var/run/nologo

xorgdir=/etc/X11
xorglog=/var/log/Xorg.0.log
xorgdrv=/usr/lib/xorg/modules/drivers

amixerrst=/usr/local/bin/reset-mixer
asoundrc=/etc/asound.conf


depend() {
    # need /etc mount
    need   localmount

    # using /var/run
    after bootmisc

    # X server is configured here
    before xdm
}


start() {
    if get_bootparam 'nologo'; then
        ebegin Disabling desktop branding
        touch ${silentflag}
        eend $?
    fi


    # Select preferred audio card based on the number of mixer controls
    if [ -e /proc/asound ]; then
        ebegin Selecting preferred audio card
        acards=`sed -n 's/^[[:blank:]]*[[:digit:]]\+[[:blank:]]*\[\([^][:blank:]]\+\)[[:blank:]]*\].*$/\1/p' /proc/asound/cards`
        bestcount=-1
        for acard in ${acards}; do
            acount=`amixer -c ${acard} controls | wc -l`
            if [ ${acount} -gt ${bestcount} ]; then
                bestcard=${acard}
                bestcount=${acount}
            fi
        done
        if [ ${bestcount} -gt -1 ]; then
            sed -i "s/\<card .*\$/card ${bestcard}/" ${asoundrc}
        else
            false
        fi
        eend $?

        ebegin Adjusting audio mixer values
        ${amixerrst}
        eend $?
    fi


    # Configure X server, unless "nox" is given to kernel
    if get_bootparam 'nox'; then
        ewarn Skipping X server configuration
    else
        ebegin "Configuring X server"
        rm -f ${xorgdir}/xorg.conf.new

        if   get_bootparam 'xvesa'; then
            tryxconf vesa        ewarn  "[fallback] VESA X server driver"
        elif get_bootparam 'xfb';   then
            tryxconf fbdev       ewarn  "[fallback] Framebuffer X server driver"
        elif get_bootparam 'xkms';  then
            tryxconf modesetting ewarn  "[fallback] Modesetting X server driver"
        else
            tryxconf native      einfo  "Native X server driver"
            tryxconf modesetting ewarn  "[fallback] Modesetting X server driver"
            tryxconf vesa        ewarn  "[fallback] VESA X server driver"
            tryxconf fbdev       ewarn  "[fallback] Framebuffer X server driver"
        fi

        # Empty last line means no section writing was stopped in the middle
        if [ -e ${xorgdir}/xorg.conf.new ]  &&  [ "`tail -1 ${xorgdir}/xorg.conf.new`" = "" ]; then
            mv ${xorgdir}/xorg.conf.new ${xorgdir}/xorg.conf

            # Turn off text cursor on vt3 (cf. xserverrc)
            echo ${BASH:+-e} -n "\033[?25l" > /dev/tty3
        else
            false
        fi
        eend $?
    fi
}


tryxconf() {
    local drvtype="$1" msgtype="$2" msg="$3"

    if [ ! -e ${xorgdir}/xorg.conf.new ]; then
        eindent
        ${msgtype} "${msg}"

        # Temporarily bind-mount specific drivers directory
        mount -rB ${xorgdrv}.${drvtype} ${xorgdrv}
        HOME=${xorgdir} Xorg -configure 2>/dev/null
        umount ${xorgdrv}

        eoutdent
    fi
}
