# src/usr/local/portage/net-misc/htpdate/files/htpdate.init
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

description="HTP time daemon"

command=/usr/sbin/htpdate
command_args="${HTPDATE_OPTS} ${PROXY} ${SERVERS}"
pidfile=/var/run/htpdate.pid

depend() {
    need net
    use dns logger
}
