# src/usr/local/portage/net-p2p/i2p/files/i2p
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

description="Starts I2P service."

pidfile=/var/run/i2p.pid
i2puser=i2p

# cf. wrapper.jvm_exit.timeout
timeout=65

wrapper=/usr/bin/wrapper
wrapper_cfg=/opt/i2p/wrapper.config


depend() {
	need net
}

start() {
	ebegin "Starting I2P"
	start-stop-daemon -S -q -u ${i2puser} -p ${pidfile} -bm -x ${wrapper} -- \
		${wrapper_cfg} wrapper.daemonize=FALSE
	eend $?
}

stop() {
	ebegin "Stopping I2P"
	start-stop-daemon -K -q -u ${i2puser} -p ${pidfile} -R ${timeout}
	eend $?
}
