# src/usr/local/portage/net-p2p/i2p/files/i2p
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

description="I2P service"

command=/usr/bin/wrapper
command_args="/opt/i2p/wrapper.config wrapper.daemonize=FALSE"
start_stop_daemon_args="-bm -u i2p"
pidfile=/var/run/i2p.pid

# cf. wrapper.jvm_exit.timeout
retry=65

depend() {
    need net
}

start_pre() {
    checkpath -q -d -m 750 -o i2p:i2p /var/run/i2p
}
